# fedora/firesim-net.service
[Unit]
Description=FireSim NIC Bringup

[Service]
Type=simple
ExecStart=/root/start-firesim-network.sh

